# pkg-plist
bin/scrutiny
bin/scrutiny-collector-metrics
etc/cron.d/scrutiny-collector
@sample %%ETCDIR%%/collector.yaml.sample
@sample %%ETCDIR%%/scrutiny.yaml.sample
%%DATADIR%%/web/3rdpartylicenses.txt
%%DATADIR%%/web/73.b0b3a038183bb5c0.js
%%DATADIR%%/web/777.100a67e5e859d9c4.js
%%DATADIR%%/web/800.2d93829b0cb17698.js
%%DATADIR%%/web/861.2916b6b81d4b6403.js
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-500.eot
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-500.svg
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-500.ttf
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-500.woff
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-500.woff2
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-600.eot
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-600.svg
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-600.ttf
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-600.woff
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-600.woff2
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-700.eot
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-700.svg
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-700.ttf
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-700.woff
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-700.woff2
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-italic.eot
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-italic.svg
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-italic.ttf
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-italic.woff
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-italic.woff2
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-regular.eot
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-regular.svg
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-regular.ttf
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-regular.woff
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono-v12-latin-regular.woff2
%%DATADIR%%/web/assets/fonts/ibm-plex-mono/ibm-plex-mono.css
%%DATADIR%%/web/assets/fonts/inter/Inter-Black.woff
%%DATADIR%%/web/assets/fonts/inter/Inter-Black.woff2
%%DATADIR%%/web/assets/fonts/inter/Inter-Bold.woff
%%DATADIR%%/web/assets/fonts/inter/Inter-Bold.woff2
%%DATADIR%%/web/assets/fonts/inter/Inter-ExtraBold.woff
%%DATADIR%%/web/assets/fonts/inter/Inter-ExtraBold.woff2
%%DATADIR%%/web/assets/fonts/inter/Inter-Italic.woff
%%DATADIR%%/web/assets/fonts/inter/Inter-Italic.woff2
%%DATADIR%%/web/assets/fonts/inter/Inter-Medium.woff
%%DATADIR%%/web/assets/fonts/inter/Inter-Medium.woff2
%%DATADIR%%/web/assets/fonts/inter/Inter-Regular.woff
%%DATADIR%%/web/assets/fonts/inter/Inter-Regular.woff2
%%DATADIR%%/web/assets/fonts/inter/Inter-SemiBold.woff
%%DATADIR%%/web/assets/fonts/inter/Inter-SemiBold.woff2
%%DATADIR%%/web/assets/fonts/inter/inter.css
%%DATADIR%%/web/assets/fonts/material-icons/MaterialIcons-Regular.codepoints
%%DATADIR%%/web/assets/fonts/material-icons/MaterialIcons-Regular.ttf
%%DATADIR%%/web/assets/fonts/material-icons/MaterialIconsOutlined-Regular.codepoints
%%DATADIR%%/web/assets/fonts/material-icons/MaterialIconsOutlined-Regular.otf
%%DATADIR%%/web/assets/fonts/material-icons/MaterialIconsRound-Regular.codepoints
%%DATADIR%%/web/assets/fonts/material-icons/MaterialIconsRound-Regular.otf
%%DATADIR%%/web/assets/fonts/material-icons/MaterialIconsSharp-Regular.codepoints
%%DATADIR%%/web/assets/fonts/material-icons/MaterialIconsSharp-Regular.otf
%%DATADIR%%/web/assets/fonts/material-icons/MaterialIconsTwoTone-Regular.codepoints
%%DATADIR%%/web/assets/fonts/material-icons/MaterialIconsTwoTone-Regular.otf
%%DATADIR%%/web/assets/fonts/material-icons/README.md
%%DATADIR%%/web/assets/fonts/material-icons/material-icons.css
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-300.eot
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-300.svg
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-300.ttf
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-300.woff
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-300.woff2
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-500.eot
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-500.svg
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-500.ttf
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-500.woff
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-500.woff2
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-700.eot
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-700.svg
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-700.ttf
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-700.woff
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-700.woff2
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-900.eot
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-900.svg
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-900.ttf
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-900.woff
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-900.woff2
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-italic.eot
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-italic.svg
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-italic.ttf
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-italic.woff
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-italic.woff2
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-regular.eot
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-regular.svg
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-regular.ttf
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-regular.woff
%%DATADIR%%/web/assets/fonts/roboto/roboto-v30-latin-regular.woff2
%%DATADIR%%/web/assets/fonts/roboto/roboto.css
%%DATADIR%%/web/assets/icons/dripicons.svg
%%DATADIR%%/web/assets/icons/feather.svg
%%DATADIR%%/web/assets/icons/heroicons-outline.svg
%%DATADIR%%/web/assets/icons/heroicons-solid.svg
%%DATADIR%%/web/assets/icons/iconsmind.svg
%%DATADIR%%/web/assets/icons/material-outline.svg
%%DATADIR%%/web/assets/icons/material-twotone.svg
%%DATADIR%%/web/assets/images/dashboard-placeholder.png
%%DATADIR%%/web/assets/images/logo/noun_Glasses_775232.svg
%%DATADIR%%/web/assets/images/logo/scrutiny-logo-dark-social.png
%%DATADIR%%/web/assets/images/logo/scrutiny-logo-dark-social.psd
%%DATADIR%%/web/assets/images/logo/scrutiny-logo-dark-text.png
%%DATADIR%%/web/assets/images/logo/scrutiny-logo-dark-text.psd
%%DATADIR%%/web/assets/images/logo/scrutiny-logo-dark.png
%%DATADIR%%/web/assets/images/logo/scrutiny-logo-dark.svg
%%DATADIR%%/web/assets/images/logo/scrutiny-logo-white-text.png
%%DATADIR%%/web/assets/images/logo/scrutiny-logo-white-text.psd
%%DATADIR%%/web/assets/images/logo/scrutiny-logo-white.psd
%%DATADIR%%/web/assets/images/logo/scrutiny-logo-white.svg
%%DATADIR%%/web/favicon-16x16.png
%%DATADIR%%/web/favicon-32x32.png
%%DATADIR%%/web/index.html
%%DATADIR%%/web/main.cf3798ce2cf0d051.js
%%DATADIR%%/web/polyfills.e5561577adeea5cb.js
%%DATADIR%%/web/runtime.6555aca5c1c1aea8.js
%%DATADIR%%/web/styles.b931f3b9dcdfc7c6.css
@dir /var/db/scrutiny
